# devops/containers/php/etc/systemd/system/shopware_consumer.target
[Install]
WantedBy=multi-user.target

[Unit]
Description=shopware_consumer service
After=mysql.service

